# Z_cron/crontabs/run_hourly
#!/system/bin/sh

for each in /etc/cron.d/hourly/* ;
do sh $each;
done;
