# Z_cron/crontabs/run_daily
#!/system/bin/sh

for each in /etc/cron.d/daily/* ;
do sh $each;
done;
