# Z_cron/crontabs/run_minutely
#!/system/bin/sh

for each in /etc/cron.d/minutely/* ;
do sh $each;
done;
